# deploy/vcenteremu.service
[Unit]
Description=vCenter Emulator (RVtools-based)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=vcenteremu
Group=vcenteremu
WorkingDirectory=/opt/vcenteremu
Environment=VCENTEREMU_PID_FILE=/var/lib/vcenteremu/run/vcenteremu.pid
Environment=VCENTEREMU_HOST=0.0.0.0
Environment=VCENTEREMU_PORT=8181
Environment=VCENTEREMU_UPLOAD_DIR=/var/lib/vcenteremu/uploads
EnvironmentFile=-/etc/vcenteremu/vcenteremu.env
ExecStartPre=/usr/bin/mkdir -p /var/lib/vcenteremu/run /var/lib/vcenteremu/uploads /var/log/vcenteremu
ExecStart=/opt/vcenteremu/.venv/bin/vcenteremu
Restart=on-failure
RestartSec=5
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
